# crontab
*/5 * * * * /ostree/build.sh >> /var/log/cron.log
# Don't remove the empty line at the end of this file. It is required to run the cron job
